# systemd/.config/systemd/user/user-programs.target

[Unit]
Description  = User Programs

[Install]
Alias = default.target
